# qy_acc.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

STOP=90
START=90

stop(){
	qy_acc.sh stop
}

start(){
	qy_acc.sh start &
}
